# ops/launchd/com.jonnycatx.major-claw.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jonnycatx.major-claw.gateway</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/pnpm</string>
    <string>--filter</string>
    <string>@majorclaw/gateway</string>
    <string>start</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Applications/Major Claw.app/Contents/Resources/workspace</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Library/Logs/MajorClaw/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Library/Logs/MajorClaw/stderr.log</string>
</dict>
</plist>

